# ge-fan.service
[Unit]
Description=Nvidia GPU Fan Controller

[Service]
ExecStart=/usr/bin/ge-fan

[Install]
WantedBy=multi-user.target
